# systemd/pipeline.timer
[Unit]
Description=Territory Pipeline Refresh Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
